# systemd/proxy-panel-api.service
[Unit]
Description=Proxy Panel API
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=proxy-panel
Group=proxy-panel
EnvironmentFile=/opt/proxy-panel/.env.generated
WorkingDirectory=/opt/proxy-panel/current
ExecStart=/opt/proxy-panel/bin/panel-api serve
Restart=always
RestartSec=2
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/proxy-panel /etc/proxy-panel /var/lib/proxy-panel /var/log/proxy-panel

[Install]
WantedBy=multi-user.target

